# debian/mopidy.postrm
#!/bin/sh
# postrm script for mopidy

set -e

if [ "$1" = "purge" ]; then
    update-rc.d mopidy remove >/dev/null
    rm -f /var/lib/mopidy/tag_cache
    rm -f /var/log/mopidy/mopidy.log
    rm -f /var/log/mopidy/mopidy-debug.log
fi

#DEBHELPER#

exit 0
